# deploy/screenguard-webui.service
[Unit]
Description=ScreenGuard Web UI
After=network-online.target screenguard-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/screenguard/webui
EnvironmentFile=/etc/screenguard/webui.env
ExecStart=/bin/sh -c 'exec /opt/screenguard/webui/venv/bin/gunicorn -w 2 -b 0.0.0.0:${UI_PORT} app:app'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
