# deploy/screenguard-server.service
[Unit]
Description=ScreenGuard Management Server
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/screenguard-server
Restart=always
RestartSec=5
WatchdogSec=60
EnvironmentFile=-/etc/screenguard/server.env

ProtectSystem=strict
ReadWritePaths=/var/lib/screenguard
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
